# scripts/crontab
# Daily snapshot at midnight Moscow time (UTC+3)
# Format: minute hour day month weekday user command
0 21 * * * root cd /workspace && DATABASE_URL="postgresql://qa_user:postgres@postgres:5432/qa_maturity" pnpm --filter @workspace/scripts run create-snapshot >> /var/log/cron-snapshot.log 2>&1
